# Server/auto_doors
#!/bin/sh /etc/rc.common
# DOORS autorun

START=50
start() {
  	#/Server/DOORS &
	/root/startDOORS.sh
} 
